# systemd/system/discache.timer
[Unit]
Description=Clean Cache

[Timer]
OnCalendar=weekly
RandomizedDelaySec=1m
Persistent=true

[Install]
WantedBy=timers.target