# systemd/system/discache.timer
[Unit]
Description=Clean Cache

[Timer]
OnCalendar=daily
RandomizedDelaySec=1m
Persistent=true