# systemd-user/dotsync.timer
[Unit]
Description=Dotfiles Synchronization

[Timer]
OnCalendar=daily
RandomizedDelaySec=1m
Persistent=true